# files/hiker-rt5350/p910nd/etc/uci-defaults/99-custom-network
#!/bin/sh
# Hiker Print - Network Configuration
# Set LAN IP to 192.168.100.1

uci set network.lan.ipaddr='192.168.100.1'
uci set network.lan.netmask='255.255.255.0'
uci commit network

